# supervisord/adblock2privoxy-supervisor.conf
[program:adblock2privoxy]
command=./opt/privoxy/privoxy-blocklist_0.2.sh
autostart = true
redirect_stderr = true
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 2
